# deploy/smartretain.service
[Unit]
Description=Smart Retain — Flask Backend
After=network.target

[Service]
User=debian
WorkingDirectory=/opt/smartretain
EnvironmentFile=/opt/smartretain/.env
ExecStart=/opt/smartretain/.venv/bin/python backend.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
